# docker/crontab
# supercronic 格式（和 crontab 几乎一样）
# 容器默认时区已设为 Asia/Shanghai（见 Dockerfile）

# 工作日多时段增量更新（净值 + 申购状态 + 股价 + 估值）
# 08:30 早间更新
30 8 * * 1-5 cd /app/scripts && python fill_missing.py && python refresh_purchase.py && python fetch_stocks.py && python calc_estimate.py && python enrich_data.py
# 17:30 午后更新
30 17 * * 1-5 cd /app/scripts && python fill_missing.py && python refresh_purchase.py && python fetch_stocks.py && python calc_estimate.py && python enrich_data.py
# 22:30 晚间更新（QDII 跨境结算 T+1，绝大多数已披露）
30 22 * * 1-5 cd /app/scripts && python fill_missing.py && python refresh_purchase.py && python fetch_stocks.py && python calc_estimate.py && python enrich_data.py

# 每月 1 日 02:00 完整流水线（含持仓、费率等慢数据）
0 2 1 * * cd /app/scripts && python scan_funds.py && python enrich_data.py && python fill_missing.py && python refresh_purchase.py && python fetch_holdings.py && python fetch_stocks.py && python calc_estimate.py
